# etc/init.d/erlnode
#! /bin/sh

# chkconfig: 2345 20 80
# description: Control an Erlang node.

exec erlnodectl "$@"
