# etc/init.d/erlnode
#! /bin/sh

# chkconfig: 2345 99 10
# description: Control an Erlang node.

exec erlnodectl "$@"
